# yt-worker/yt-worker.service
[Unit]
Description=Olymp/YT Worker service

[Service]
Type=simple
ExecStart=/usr/bin/yt-worker
StandardError=journal
