# worker/crontab
*/2 * * * * refresh_auth_details

29 21 * * * create_lesson_payments

* * * * * frontier_schedule
8 * * * * * discover_csts_athletes
